# launchd/com.recalldeck.clawd.orchestrator-guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.recalldeck.clawd.orchestrator-guardian</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>/Users/joshwegener/clawd/scripts/run_orchestrator_guardian.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>60</integer>

    <key>WorkingDirectory</key>
    <string>/Users/joshwegener/clawd</string>

    <key>StandardOutPath</key>
    <string>/Users/joshwegener/clawd/memory/launchd-orchestrator-guardian.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/joshwegener/clawd/memory/launchd-orchestrator-guardian.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
  </dict>
</plist>

